# deploy/rag-service.service
[Unit]
Description=RAG Service Web UI + Background Ingestion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/egarcia/.rag-service/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/egarcia/.rag-service/app/server.py
Restart=always
RestartSec=3
KillMode=process
KillSignal=SIGINT
FinalKillSignal=SIGKILL
SendSIGKILL=yes
TimeoutStopSec=5

[Install]
WantedBy=default.target
