# deploy/rag-service.service
[Unit]
Description=RAG Service Web UI + Background Ingestion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/egarcia/.rag-service/app
Environment=PYTHONUNBUFFERED=1
Environment=RAG_SERVICE_HOME=/home/egarcia/.rag-service/data
Environment=RAG_SERVICE_HOST=0.0.0.0
Environment=RAG_SERVICE_PORT=6334
ExecStart=/usr/bin/python3 /home/egarcia/.rag-service/app/server.py
Restart=always
RestartSec=3
KillMode=control-group
KillSignal=SIGINT
FinalKillSignal=SIGKILL
SendSIGKILL=yes
TimeoutStopSec=10

[Install]
WantedBy=default.target
